# .config/systemd/user/check_updates.service
[Unit]
Description=Check for system updates

[Service]
ExecStart=%h/.scripts/check_for_updates.sh
